# extra/zerotier/files/zerotier-one.service
[Unit]
Description=ZeroTier One

[Service]
ExecStart=/usr/bin/zerotier-one
Restart=always
KillMode=process

[Install]
WantedBy=multi-user.target
