# core/hokuto/files/hokuto-bump.timer
[Unit]
Description=Run hokuto bump every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
